# dbdump.cron
# scheduling suggestion, everyday at 2AM
0 2 * * *       root    bash /usr/local/bin/dbdump/dbdump
